# debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove|purge)
        systemctl stop ergo-proxy-node 2>/dev/null || true
        systemctl disable ergo-proxy-node 2>/dev/null || true
        ;;
esac
